# pkgbuilds/tsw-printer/tsw-printer.install
post_package() {
    systemctl enable cups.socket
}
