# pkgbuilds/tsw-base/tsw-base.install
post_package() {
    systemctl enable reflector.timer
    systemctl enable paccache.timer
}
