# pkgbuilds/tsw-base/tsw-base.install
post_package() {
    systemctl enable reflector.time
}
