# pkgbuilds/tsw-network/tsw-network.install
# Maintainer: Your Name <youremail@domain.com>

_service="NetworkManager"

post_install=(
    sudo systemctl enable $_service
    sudo systemctl start $_service
)

pre_remove=(
    sudo systemctl disable $_service
    sudo systemctl stop $_service
)
