# pkgbuilds/tsw-desktop/tsw-desktop.install
post_package() {
    systemctl enable sddm.service
}
